# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:passion]
command=go run main.go
autostart=true
autorestart=true
stopsignal=TERM
stdout_logfile=log/bot.stdout.log
stderr_logfile=log/bot.stderr.log
